# systemd/nanoclaw@.service
[Unit]
Description=NanoClaw Agent (%i)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/home/pineappleles/nanoclaw
EnvironmentFile=/home/pineappleles/nanoclaw/systemd/env-%i.conf
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nanoclaw-%i

[Install]
WantedBy=default.target
